# 99-reboot-usb-devices.rules
SUBSYSTEM=="usb", ACTION=="remove", ATTRS{idVendor}==04b0", ATTRS{idProduct}==”044c", RUN+="bash /opt/reboot-usb-devices/trigger.sh"
